# cron.hourly/rpisc-cleanup
#!/bin/sh
sudo -u motion /opt/rpi_security_cam/sendemail.py -c /motion
